# var/lib/dpkg/info/fuse-2.x-x64-dev.list
./
./usr/
./usr/local/
./usr/local/x86_64-pc-linux-gnu/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/sbin/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/sbin/mount.fuse
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/etc/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/etc/init.d/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/etc/init.d/fuse
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/etc/udev/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/etc/udev/rules.d/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/etc/udev/rules.d/99-fuse.rules
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/share/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/share/man/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/share/man/man1/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/share/man/man1/fusermount.1
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/share/man/man1/ulockmgr_server.1
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/share/man/man8/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/share/man/man8/mount.fuse.8
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/lib/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/lib/libulockmgr.so.1.0.1
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/lib/libfuse.so
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/lib/pkgconfig/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/lib/pkgconfig/fuse.pc
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/lib/libfuse.so.2
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/lib/libulockmgr.so
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/lib/libulockmgr.so.1
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/lib/libfuse.so.2.9.4
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/include/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/include/ulockmgr.h
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/include/fuse.h
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/include/fuse/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/include/fuse/fuse.h
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/include/fuse/fuse_compat.h
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/include/fuse/fuse_lowlevel_compat.h
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/include/fuse/cuse_lowlevel.h
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/include/fuse/fuse_opt.h
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/include/fuse/fuse_lowlevel.h
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/include/fuse/fuse_common.h
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/include/fuse/fuse_common_compat.h
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/bin/
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/bin/fusermount
./usr/local/x86_64-pc-linux-gnu/x86_64-pc-linux-gnu/sys-root/usr/bin/ulockmgr_server
